# i3/.config/systemd/user/newsboat.service
[Unit]
Description=Check RSS feeds

[Service]
Type=oneshot
ExecStart=%h/.scripts/newsboat_update
StandardOutput=journal
