# i3/.config/systemd/user/pacupdate.service
[Unit]
Description=Check number of available packages

[Service]
Type=oneshot
ExecStart=%h/.scripts/pacman_update
StandardOutput=journal
