# old/configs/config/systemd/user/pacupdate.timer
[Unit]
Description=Check number of available packages

[Timer]
OnBootSec=15min
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
