# systemd/updatedb.service
[Unit]
Description=Update file locate index

[Service]
Type=oneshot
ExecStart=updatedb
StandardOutput=journal
