# systemd/updatedb.timer
[Unit]
Description=Update file locate index

[Timer]
OnBootSec=15min
OnUnitActiveSec=30min

[Install]
WantedBy=timers.target
