# configs/config/systemd/user/home_backup.timer
[Unit]
Description=Backup home directory to OneDrive folder

[Timer]
OnCalendar=*-*-* 1:00:00

[Install]
WantedBy=timers.target
